# deploy/systemd/dakwah-celery-worker@.service
[Unit]
# Templated unit — runs the Celery worker (ingest + classification + cluster).
# Concurrency=2 is enough at prototype scale; bump if the ingest queue
# backs up. Each env has its own worker so dev can't poison prod's queue.
Description=Dakwah-Lens Celery worker — %i
After=network.target docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/srv/dakwah-lens-%i/api
EnvironmentFile=/srv/dakwah-lens-%i/.env
# Cap PyTorch / MKL / OpenBLAS thread pools so IndoBERT doesn't starve
# the rest of the worker (Postgres, Celery, the API server) on the shared
# 4-vCPU VM. Setting these here means they're in place BEFORE Python
# imports torch — the most reliable point.
# Currently 1 — IndoBERT runs single-threaded; bump only after measuring
# that ingest latency is actually the bottleneck.
Environment="OMP_NUM_THREADS=1"
Environment="MKL_NUM_THREADS=1"
Environment="OPENBLAS_NUM_THREADS=1"
Environment="TOKENIZERS_PARALLELISM=false"
ExecStart=/home/deploy/.local/bin/uv run celery -A api.workers.celery_app worker --loglevel=info --concurrency=2
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
