# deploy/systemd/dakwah-celery-beat@.service
[Unit]
# Templated unit — runs the Celery beat scheduler that triggers ingest +
# recluster ticks. Exactly ONE beat process per env is required, otherwise
# every scheduled task fires multiple times. systemd makes that easy: one
# instance per env, no scaling.
Description=Dakwah-Lens Celery beat — %i
After=network.target docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/srv/dakwah-lens-%i/api
EnvironmentFile=/srv/dakwah-lens-%i/.env
ExecStart=/home/deploy/.local/bin/uv run celery -A api.workers.celery_app beat --loglevel=info
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
